# examples/systemd/cm3500-capacity.path
[Unit]
Description=Watch CM3500 capacity output for shaping automation

[Path]
PathExists=/run/cm3500/capacity.json
PathChanged=/run/cm3500/capacity.json
Unit=cm3500-capacity.service

[Install]
WantedBy=multi-user.target
